# supervisord.conf
[supervisord]
nodaemon=true
loglevel=error

[program:chronos]
command=chronos /etc/chronos.cfg
user=chronos
autostart=true
autorestart=true
